# app/protos/crontab
@reboot /bin/bash @CONFIG@/start.sh >> @LOGS@/reboot.log 2>&1 &
*/1 * * * * /bin/bash @CONFIG@/healthcheck.sh >> @LOGS@/healthcheck.log 2>&1 &
